# an_go.conf
description "Service for running an_go scheduler and filewatcher on startup/reboot"
author "Will and Andy"

start on filesystem or runlevel [2345]
stop on shutdown
respawn

script

    export HOME="/home/ubuntu/an_go"
    echo $$ > /var/run/an_go.pid
    cd /home/ubuntu/an_go
    exec nohup node /home/ubuntu/an_go/an_go.js

end script

pre-start script
    echo "[`date`] an_go Starting" >> /var/log/an_go.log
end script

pre-stop script
    rm /var/run/an_go.pid
    echo "[`date`] an_go Stopping" >> /var/log/an_go.log
end script